# root/etc/hotplug.d/ieee80211/10-wifi-detect
#!/bin/sh

# just skip wifi config in openwrt of rtd1296
exit 0

[ "${ACTION}" = "add" ] && {
	/sbin/wifi config
}
