# root/usr/lib/opkg/info/aria2.list
/etc/config/aria2
/etc/init.d/aria2
/usr/bin/aria2c
