# root/usr/lib/opkg/info/avahi-dbus-daemon.list
/usr/sbin/avahi-daemon
/etc/init.d/avahi-daemon
/etc/avahi/avahi-daemon.conf
