# root/usr/lib/opkg/info/forked-daapd.list
/usr/lib/forked-daapd/forked-daapd-sqlext.la
/usr/sbin/forked-daapd
/usr/lib/forked-daapd/forked-daapd-sqlext.so
/etc/forked-daapd.conf
/etc/init.d/forked-daapd
