# root/usr/lib/opkg/info/logd.list
/sbin/logread
/sbin/logd
/etc/init.d/log
