# root/usr/lib/opkg/info/luci-base.list
/www/luci-static/resources/cbi/user.gif
/www/luci-static/resources/firewall.js
/usr/lib/lua/luci/view/view.htm
/usr/lib/lua/luci/store.lua
/etc/config/luci
/www/luci-static/resources/rpc.js
/www/luci-static/resources/icons/alias.png
/www/luci-static/resources/cbi/remove.gif
/www/luci-static/resources/form.js
/www/luci-static/resources/icons/alias_disabled.png
/www/luci-static/resources/icons/switch_disabled.png
/usr/lib/lua/luci/template.lua
/etc/init.d/ucitrack
/www/luci-static/resources/icons/loading.gif
/usr/lib/lua/luci/view/header.htm
/usr/lib/lua/luci/view/empty_node_placeholder.htm
/www/luci-static/resources/icons/signal-75-100.png
/www/luci-static/resources/icons/signal-none.png
/sbin/luci-reload
/usr/lib/lua/luci/view/footer.htm
/usr/lib/lua/luci/controller/admin/uci.lua
/www/luci-static/resources/ui.js
/www/luci-static/resources/icons/wifi.png
/usr/lib/lua/luci/view/error404.htm
/www/luci-static/resources/cbi/save.gif
/usr/lib/lua/luci/debug.lua
/usr/lib/lua/luci/sys/zoneinfo.lua
/www/luci-static/resources/network.js
/www/luci-static/resources/tools/prng.js
/www/luci-static/resources/fs.js
/www/luci-static/resources/icons/port_up.png
/usr/lib/lua/luci/config.lua
/usr/lib/lua/luci/util.lua
/usr/lib/lua/luci/sys.lua
/usr/lib/lua/luci/view/csrftoken.htm
/www/luci-static/resources/cbi/file.gif
/www/luci-static/resources/protocol/dhcp.js
/etc/config/ucitrack
/www/luci-static/resources/promis.min.js
/usr/lib/lua/luci/sys/zoneinfo/tzoffset.lua
/www/luci-static/resources/cbi/download.gif
/www/luci-static/resources/icons/port_down.png
/www/luci-static/resources/cbi/reset.gif
/usr/lib/lua/luci/sys/zoneinfo/tzdata.lua
/usr/share/acl.d/luci-base.json
/www/luci-static/resources/cbi/help.gif
/www/luci-static/resources/icons/tunnel.png
/www/luci-static/resources/xhr.js
/www/luci-static/resources/icons/vlan_disabled.png
/www/luci-static/resources/icons/vlan.png
/www/luci-static/resources/cbi/reload.gif
/usr/lib/lua/luci/controller/admin/index.lua
/lib/upgrade/keep.d/luci-base
/www/luci-static/resources/uci.js
/usr/lib/lua/luci/ltn12.lua
/usr/lib/lua/luci/ccache.lua
/www/luci-static/resources/cbi.js
/usr/lib/lua/luci/model/uci.lua
/www/luci-static/resources/cbi/find.gif
/usr/lib/lua/luci/cacheloader.lua
/www/cgi-bin/luci
/usr/lib/lua/luci/tools/webadmin.lua
/www/luci-static/resources/icons/signal-50-75.png
/www/luci-static/resources/cbi/edit.gif
/etc/luci-uploads/.placeholder
/usr/lib/lua/luci/i18n.lua
/usr/lib/lua/luci/version.lua
/www/luci-static/resources/tools/widgets.js
/www/luci-static/resources/protocol/none.js
/usr/lib/lua/luci/view/sysauth.htm
/usr/lib/lua/luci/dispatcher.lua
/www/luci-static/resources/icons/wifi_disabled.png
/www/luci-static/resources/icons/signal-0-25.png
/www/luci-static/resources/cbi/fieldadd.gif
/usr/lib/lua/luci/http.lua
/www/luci-static/resources/cbi/apply.gif
/www/luci-static/resources/cbi/key.gif
/www/luci-static/resources/cbi/link.gif
/usr/libexec/rpcd/luci
/www/luci-static/resources/cbi/arrow.gif
/usr/share/rpcd/acl.d/luci-base.json
/usr/lib/lua/luci/tools/proto.lua
/www/luci-static/resources/protocol/static.js
/www/luci-static/resources/validation.js
/usr/lib/lua/luci/sgi/uhttpd.lua
/www/luci-static/resources/icons/tunnel_disabled.png
/usr/lib/lua/luci/template/parser.so
/www/luci-static/resources/icons/signal-0.png
/www/luci-static/resources/icons/switch.png
/www/luci-static/resources/icons/ethernet_disabled.png
/www/luci-static/resources/cbi/down.gif
/www/luci-static/resources/cbi/folder.gif
/www/luci-static/resources/icons/bridge_disabled.png
/www/luci-static/resources/luci.js
/www/luci-static/resources/cbi/up.gif
/usr/lib/lua/luci/view/error500.htm
/www/luci-static/resources/icons/ethernet.png
/www/index.html
/www/luci-static/resources/icons/signal-25-50.png
/usr/lib/lua/luci/view/indexer.htm
/usr/lib/lua/luci/sgi/cgi.lua
/www/luci-static/resources/cbi/add.gif
/www/luci-static/resources/icons/bridge.png
