# root/usr/lib/opkg/info/mdadm.list
/etc/init.d/mdadm
/sbin/mdadm
/etc/config/mdadm
