# root/usr/lib/opkg/info/minidlna.list
/etc/init.d/minidlna
/usr/bin/minidlna
/etc/config/minidlna
