# root/usr/lib/opkg/info/netatalk.list
/usr/lib/uams/uams_clrtxt.so
/usr/lib/uams/uams_dhx.so
/etc/extmap.conf
/usr/lib/libatalk.so.18.0.0
/lib/upgrade/keep.d/netatalk
/etc/init.d/afpd
/usr/lib/uams/uams_passwd.so
/usr/bin/dbd
/usr/sbin/cnid_metad
/usr/lib/uams/uams_randnum.so
/usr/bin/afppasswd
/usr/lib/uams/uams_dhx2_passwd.so
/usr/lib/libatalk.so.18
/usr/bin/ad
/usr/lib/uams/uams_guest.so
/usr/lib/uams/uams_dhx_passwd.so
/usr/lib/libatalk.so
/usr/sbin/cnid_dbd
/usr/lib/uams/uams_dhx2.so
/etc/afp.conf
/usr/sbin/afpd
