# root/usr/lib/opkg/info/nfs-kernel-server.list
/etc/exports
/usr/sbin/exportfs
/etc/init.d/nfsd
/usr/sbin/rpc.nfsd
/usr/sbin/rpc.mountd
/usr/sbin/rpc.statd
/usr/sbin/sm-notify
