# root/usr/lib/opkg/info/rpcbind.list
/usr/sbin/rpcbind
/usr/bin/rpcinfo
/etc/init.d/rpcbind
