# root/usr/lib/opkg/info/transmission-daemon-mbedtls.list
/etc/init.d/transmission
/etc/sysctl.d/20-transmission.conf
/usr/bin/transmission-daemon
/etc/config/transmission
