# root/usr/lib/opkg/info/urandom-seed.list
/lib/preinit/81_urandom_seed
/etc/init.d/urandom_seed
/sbin/urandom_seed
