# root/usr/lib/opkg/info/vsftpd.list
/etc/vsftpd.conf
/lib/upgrade/keep.d/vsftpd
/usr/sbin/vsftpd
/etc/init.d/vsftpd
