# root/etc/uci-defaults/10-fstab
[ ! -f /etc/config/fstab ] && ( /sbin/gen_fstab > /dev/null )
exit 0
